# docker/cron/crontab
59 23 * * * root /srv/app/bin/script.sh
# An empty line